# data/ssserv.ini
[program:shadowsocks]
command=ssserver -c /etc/supervisor.d/ssserv.conf
autorestart=true
user=nobody
